# build/debian/postrm
#!/bin/sh

#
# deinstall optional files
#
rcdir=/etc/init.d
logrotatedir=/etc/logrotate.d

if [ -n "`which chkconfig`" ]
then
	chkconfig --del nbspd
elif [ -n "`which update-rc.d`" ]
then
	update-rc.d -f nbspd remove
fi

cat <<__EOF__
To deinstall completely this package, the noaaport user and the directories

/usr/local/etc/nbsp
/usr/local/share/nbsp
/usr/local/libexec/nbsp
/var/noaaport
/var/log/noaaport
/var/log/nbsp
/var/run/nbsp

must be removed manually.
__EOF__
